# client/systemd/ollama-tunnel.service
[Unit]
Description=SSH tunnel from laptop localhost:11434 to workstation localhost:11434
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=LOCAL_OLLAMA_BIND_HOST=127.0.0.1
EnvironmentFile=%h/.config/local-agent-devstack/client.env
ExecStart=/usr/bin/ssh -N -o ExitOnForwardFailure=yes -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -L ${LOCAL_OLLAMA_BIND_HOST}:${LOCAL_OLLAMA_PORT}:localhost:${OLLAMA_PORT} ${WORKSTATION_USER}@${WORKSTATION_HOST}
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
